# usbip/etc/services.d/usbip/run
#!/command/with-contenv bashio

bashio::log.info "attaching usbip devices"

for device in $(bashio::config 'devices|keys'); do
	server_address=$(bashio::config "devices[${device}].server_address")
	bus_id=$(bashio::config "devices[${device}].bus_id")

	bashio::log.info "attaching ${server_address} ${bus_id}"

	/usr/sbin/usbip --debug attach -r ${server_address} -b ${bus_id}
done

# wait forever
exec tail -f /dev/null
